# deploy/systemd/lake-dispatch@.service
# Runs every enabled source for one schedule. %i is daily|weekly|monthly|yearly.
# The timers activate this, not lake-scrape@ directly, so adding a source is a
# YAML edit and `lake sync-sources` — no new unit, no `systemctl enable`.

[Unit]
Description=lake dispatcher: %i schedule
Documentation=file:/opt/lake/docs/runbook.md
After=network-online.target mnt-nas.mount postgresql.service
Wants=network-online.target
Requires=mnt-nas.mount
OnFailure=lake-alert@dispatch-%i.service

[Service]
Type=oneshot
User=lake
Group=lake
WorkingDirectory=/opt/lake
EnvironmentFile=/etc/lake/lake.env
ExecStart=/opt/lake/.venv/bin/lake scrape-schedule %i
# Yearly dumps can take hours. Daily jobs finish in seconds.
TimeoutStartSec=21600

Nice=10
IOSchedulingClass=idle
MemoryMax=4G

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
ReadWritePaths=/mnt/nas/lake /var/log/lake /var/lib/lake
