# deploy/systemd/lake-retry.timer
[Unit]
Description=lake: retry failed runs every 30 minutes

[Timer]
OnCalendar=*:0/30
RandomizedDelaySec=120
Persistent=false
Unit=lake-retry.service

[Install]
WantedBy=timers.target
