# deploy/systemd/lake-retry.service
# Cross-run retry. Picks up runs that failed, are under the attempt cap, and
# have not since succeeded. Each retry creates a NEW run row with attempt=N+1,
# so the whole history stays visible in the catalog.

[Unit]
Description=lake: retry failed runs
After=network-online.target mnt-nas.mount postgresql.service
Requires=mnt-nas.mount

[Service]
Type=oneshot
User=lake
Group=lake
WorkingDirectory=/opt/lake
EnvironmentFile=/etc/lake/lake.env
ExecStart=/opt/lake/.venv/bin/lake retry --max-attempts 3 --older-than-minutes 15
TimeoutStartSec=3600

Nice=15
IOSchedulingClass=idle
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/mnt/nas/lake /var/log/lake /var/lib/lake
